# deploy/vlm-db.service
[Unit]
Description=VLM DB Server (SQLite Event Storage)
After=network.target

[Service]
Type=simple
User=vlmapp
Group=vlmapp
WorkingDirectory=/home/ubuntu/hio_intelligence_stream
Environment=PYTHONPATH=/home/ubuntu/hio_intelligence_stream
EnvironmentFile=/home/ubuntu/hio_intelligence_stream/.env

ExecStart=/home/ubuntu/hio_intelligence_stream/venv/bin/uvicorn \
    db_server.main:app \
    --host 127.0.0.1 \
    --port 8001 \
    --workers 1 \
    --timeout-keep-alive 65 \
    --log-level info

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vlm-db

[Install]
WantedBy=multi-user.target
